# supervisord.conf
# TrendXL 2.0 - Supervisor Configuration for Fullstack Application
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=debug
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
user=root
silent=false

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700
username=dummy
password=dummy

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock
username=dummy
password=dummy

# Python FastAPI Backend Service
[program:backend]
command=/app/venv/bin/python run_server.py
directory=/app/backend
user=trendxl
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/backend.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
startsecs=5
startretries=3
environment=PYTHONPATH="/app/backend",PYTHONUNBUFFERED="1",HOST="127.0.0.1",PORT="8000",VIRTUAL_ENV="/app/venv",PATH="/app/venv/bin:/usr/local/bin:/usr/bin:/bin"
priority=100
stderr_logfile=/var/log/supervisor/backend_error.log
stderr_logfile_maxbytes=50MB

# Nginx Frontend Service  
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/nginx.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
user=root
startsecs=2
startretries=3
priority=200
stderr_logfile=/var/log/supervisor/nginx_error.log
stderr_logfile_maxbytes=50MB
